# devops/ansible/.service
[Unit]
Description=Run go server in background


[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/rudeshell/rudeshell_go
ExecStart=/home/ubuntu/rudeshell/rudeshell_go/run.sh
